# deploy/ip-geo-exporter.service
[Unit]
Description=ip-geo-exporter - Per-IP traffic monitoring with GeoIP
Documentation=https://github.com/kof8855/ip-geo-exporter
After=network.target nftables.service

[Service]
Type=simple
WorkingDirectory=/opt/ip-geo-exporter
ExecStart=/usr/local/bin/ip-geo-exporter --config=/opt/ip-geo-exporter/config.yaml

# 如果不想用 config 文件，也可以直接传参：
# ExecStart=/usr/local/bin/ip-geo-exporter \
#     --listen-address=:9100 \
#     --poll-interval=5s \
#     --protocol-filter=tcp,udp,icmp \
#     --enable-ipv6=true \
#     --filter-private=false \
#     --geoip-lang=zh-CN \
#     --cleanup-on-exit=true \
#     --log-level=info

Restart=always
RestartSec=5

# nftables 需要 CAP_NET_ADMIN 来操作规则集
CapabilityBoundingSet=CAP_NET_ADMIN
AmbientCapabilities=CAP_NET_ADMIN

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/run
PrivateTmp=true

[Install]
WantedBy=multi-user.target
